# samples/scripts/mha-service
#! /bin/sh
#
# $Id:$
#
# mha:       Starts the masterha_manager as Daemon
#
# Version:     1.0 
#
# chkconfig: 2345 90 10
# description: Start & Ends masterha_manager as a Daemon \
#
# processname: masterha_manager
# config: /etc/mha/

. /lib/lsb/init-functions

CLUSTER='group';
CONF=/etc/mha/${CLUSTER}.cnf  
LOGFILE=/var/log/masterha/${CLUSTER}/${CLUSTER}.log  

start () {
    echo "Starting $CLUSTER..."
    command_ls=`nohup masterha_manager --conf=$CONF --manager_log=${LOGFILE} < /dev/null >> $LOGFILE 2>&1 &`
    if [[ $? eq 0 ]]
    then 
        echo 'Success';
    else
        echo 'Failed';
    fi
}
stop () {
    echo "Stopping $CLUSTER"
    start-stop-daemon --stop --quiet --oknodo --pidfile "$pid_file"
}
status () {
    echo "Status of $CLUSTER"
    masterha_check_status --conf=${CONF}

}
case $1 in
start)
    if status; then exit 0; fi
    start
    status && exit 0 || exit $?
    ;;
    stop)
    stop
    ;;
    reload)
    stop
    start
    ;;
    restart)
    stop
    start
    ;;
    status)
    status && exit 0 || exit $?
    ;;
    *)
    echo "Usage: $0 {start|stop|restart|reload|status}"
    exit 1
    ;;
    esac
    exit 0
